# playbooks/files/docker-stats.timer
[Unit]
Description=Continuous Docker stats monitoring

[Timer]
OnUnitInactiveSec=5
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
