# playbooks/files/docker-stats.service
[Unit]
Description=Docker stats monitoring

[Service]
Type=simple
ExecStart=/bin/docker stats
RemainAfterExit=no
Restart=no
