# playbooks/files/mpstat.service
[Unit]
Description=CPU monitoring with mpstat

[Service]
Type=simple
ExecStart=/bin/mpstat 1
RemainAfterExit=no
Restart=no
