# playbooks/files/vmstat.service
[Unit]
Description=System monitoring with vmstat

[Service]
Type=simple
ExecStart=/bin/vmstat --unit M --one-header 1
RemainAfterExit=no
Restart=no
